# solar-dashboard.service
[Unit]
Description=Solar PV E-Ink Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/solar-dashboard
ExecStart=/usr/bin/python3 /home/pi/solar-dashboard/solar_dashboard.py
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal

# Variabili ambiente (opzionale)
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
